# scripts/amentia.conf
[program:amentia]
command = /home/andres/projects/amentia/scripts/start_gunicorn.sh
user = ubuntu
stdout_logfile = /home/andres/projects/amentia/logs/gunicorn_supervisor.log
redirect_stderr = true